# rootfs/etc/periodic/daily/fetch-dat
#!/bin/sh

cd /var/www
php scripts/fetch-dat.php --mode recent
