# rootfs/etc/service/2chproxy/run
#!/bin/sh

rm -f /tmp/2chproxy.pid
exec /usr/bin/sudo -u nobody /usr/local/bin/2chproxy.pl -c /etc/2chproxy.conf
